# base/xinetd.d/container
# default: on
# description: Basic container controller
service container
{
   port            = 10101
   socket_type     = stream
   wait            = no
   user            = root
   server          = /container/controller
   disable         = no
}

